# systemd/whatsapp-bot.service
[Unit]
Description=WhatsApp Bot - AI-powered message responder
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/whatsapp-bot
ExecStart=/path/to/whatsapp-bot/run.sh
Restart=always
RestartSec=10
StandardOutput=append:/path/to/whatsapp-bot/logs/service.log
StandardError=append:/path/to/whatsapp-bot/logs/service.log

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
